# server/api-scheduler/deploy/supervisord/api-scheduler.ini
[program:api-scheduler]
directory = /home/ec2-user/sd-async-inference/server/api-scheduler
command =  python3 main.py
autostart = true
startsecs = 5
autorestart =true
startretries = 3
user = ec2-user 
redirect_stderr = false
stdout_logfile_maxbytes =20MB
stdout_logfile_backups =20
stdout_logfile=/home/ec2-user/sd-async-inference/server/api-scheduler/deploy/log/api-schedluer.log 
stderr_logfile=/home/ec2-user/sd-async-inference/server/api-scheduler/deploy/log/api-schedluer.err.log 