# server/api-scheduler/deploy/supervisord/stablediffusion-webui.ini
[program:sd-webui]
directory =/home/ec2-user/stable-diffusion-webui
command =  bash webui.sh --api --listen --enable-insecure-extension-access
autostart = true
startsecs = 5
autorestart =true
startretries = 3
user = ec2-user 
redirect_stderr = false
stdout_logfile_maxbytes =20MB
stdout_logfile_backups =20
stdout_logfile=/home/ec2-user/sd-async-inference/server/api-scheduler/deploy/log/sd-webui.log 
stderr_logfile=/home/ec2-user/sd-async-inference/server/api-schedluer/deploy/log/sd-webui.err.log 